# launchd/com.yazydzhi.meeting-automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yazydzhi.meeting-automation</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/azg/repository/meeting_automation/venv/bin/python</string>
        <string>/Users/azg/repository/meeting_automation/src/service_manager.py</string>
        <string>--interval</string>
        <string>300</string>
        <string>--media-interval</string>
        <string>1800</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/azg/repository/meeting_automation</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/azg/repository/meeting_automation/venv/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/azg/repository/meeting_automation/logs/service.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/azg/repository/meeting_automation/logs/service_error.log</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>10</integer>
    
    <key>StartInterval</key>
    <integer>300</integer>
</dict>
</plist>
